# bin/home-deck-games-bb.lan.automount
[Unit]
Description=Automount Steam NFS Share

[Automount]
Where=/home/deck/games/bb.lan/

[Install]
WantedBy=multi-user.target

